# deployment/care-calendar.service
[Unit]
Description=Care Calendar Flask App
After=network.target

[Service]
User=monroe
WorkingDirectory=/home/monroe/Care-Calendar
Environment=HOST=0.0.0.0
Environment=PORT=8080
Environment=FLASK_DEBUG=0
# Set a strong secret in production
Environment=FLASK_SECRET_KEY=change-me-to-a-long-random-string
# Relocated DB path (Phase B). Ensure file exists at this path before restarting.
Environment=CARE_DB_PATH=/home/monroe/Care-Calendar/data/database.db
ExecStart=/home/monroe/Care-Calendar/.venv/bin/python /home/monroe/Care-Calendar/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
